# systemd/meshconfig.service
[Unit]
Description=MeshConfig Workbench static web server
After=network.target

[Service]
Type=simple
DynamicUser=yes
WorkingDirectory=/opt/meshconfig
Environment=HOST=0.0.0.0
Environment=PORT=8420
ExecStart=/opt/meshconfig/scripts/launch.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
